# agent/deploy/homelytics-agent.service
[Unit]
Description=Homelytics Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/homelytics-agent -config /etc/homelytics/agent.yaml
Restart=always
RestartSec=5s
User=root
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/homelytics /var/log/homelytics

[Install]
WantedBy=multi-user.target
